# deploy/systemd/secumon-ingestion.service
[Unit]
Description=SecuMon Ingestion Service
Documentation=https://github.com/secuaas/secumon
After=network.target timescaledb.service postgresql.service

[Service]
Type=simple
User=secumon
Group=secumon
WorkingDirectory=/opt/secumon

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/secumon/ingestion.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/secumon

# Execution
ExecStart=/usr/local/bin/secumon-ingestion

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secumon-ingestion

[Install]
WantedBy=multi-user.target
